# templates/default/init_system/systemd/airflow-kerberos.service.erb
[Unit]
Description=Airflow kerberos daemon
After=network.target

[Service]
EnvironmentFile=/etc/sysconfig/airflow
User=<%= @user %>
Group=<%= @group %>
ExecStart=<%= @bin_path %>/airflow kerberos --pid <%= @run_path %>/airflow-kerberos.pid
Restart=always

[Install]
WantedBy=multi-user.target