# templates/default/init_system/systemd/airflow-scheduler.service.erb
[Unit]
Description=Airflow scheduler daemon
After=network.target

[Service]
EnvironmentFile=/etc/sysconfig/airflow
User=<%= @user %>
Group=<%= @group %>
ExecStart=<%= @bin_path %>/airflow scheduler --pid <%= @run_path %>/airflow-scheduler.pid
Restart=always

[Install]
WantedBy=multi-user.target