# templates/default/init_system/systemd/airflow-webserver.service.erb
[Unit]
Description=Airflow webserver daemon
After=network.target

[Service]
EnvironmentFile=/etc/sysconfig/airflow
User=<%= @user %>
Group=<%= @group %>
ExecStart=<%= @bin_path %>/airflow webserver --pid <%= @run_path %>/airflow-webserver.pid
Restart=always

[Install]
WantedBy=multi-user.target